# createdb.conf
[program:createdb]
command=/createdb.sh
autostart=true
autorestart=false
stderr_logfile=/var/log/createdb.err.log
stdout_logfile=/var/log/createdb.out.log